# musikbox.service
[Unit]
Description=musikbox – folder-based music player
After=sound.target

[Service]
Type=simple
ExecStart=/usr/local/bin/musikbox
Restart=on-failure
User=musikbox

[Install]
WantedBy=multi-user.target
